# clients_server_automation/conductor/qabu-conductor.service
[Unit]
Description=Qabu conductor — client lifecycle manager
After=docker.service
Requires=docker.service

[Service]
Type=simple
User=brande
ExecStart=/home/brande/app/conductor
RuntimeDirectory=qabu
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
